# plugins/wheezy/sshkeys/files/etc/network/if-up.d/zz-ssh_keys
#!/bin/sh

set -e

if ! ls /etc/ssh/ssh_host* 1> /dev/null 2>&1; then
	/usr/sbin/dpkg-reconfigure -f noninteractive openssh-server
	/etc/init.d/ssh restart
fi
