# files/common/etc/network/if-up.d/zz-ssh_keys
#!/bin/bash

set -e

if [ ! -e /etc/ssh/ssh_host* ]; then
	/usr/sbin/dpkg-reconfigure -f noninteractive openssh-server
	/etc/init.d/ssh restart
fi
